# src/settings/supervisor.conf
[program:dinner]
command=gunicorn wsgi --bind=0.0.0.0:9002
directory=/home/django/dinner/src
daemon=true
workers=2
autorestart=true
user=nobody